# roles/zfs/templates/zfs-mount-encrypted-datasets.service.j2
[Unit]
Description=Mount encrypted ZFS datasets
Documentation=man:zfs(8)
DefaultDependencies=no
After=systemd-udev-settle.service
After=zfs-import-cache.service
After=zfs-load-keys.service
After=systemd-remount-fs.service
Before=local-fs.target
ConditionPathIsDirectory=/sys/module/zfs

[Service]
Type=oneshot
RemainAfterExit=no
ExecStartPre=/bin/sleep 3s
ExecStart=/sbin/zfs mount -a -v

[Install]
WantedBy=zfs.target