# roles/restic/templates/10-override.timer.j2
[Timer]
OnUnitActiveSec={{ item.1.systemd_timer_interval }}
