# roles/restic/templates/restic@.timer.j2
[Unit]
Description=Restic %i backup timer

[Timer]
OnStartupSec=5m
Unit=restic@%i.service

[Install]
WantedBy=timers.target
