# roles/gad/templates/10-override.timer.j2
[Timer]
OnUnitActiveSec={{ gad_systemd_timer_interval }}