# roles/gad/templates/gad6.timer.j2
[Unit]
Description=Update Gandi DNS Records

[Timer]
OnStartupSec=5m
OnUnitActiveSec=15m
Unit=gad6.service

[Install]
WantedBy=timers.target