# roles/gad/templates/gad@.timer.j2
[Unit]
Description=Update Gandi DNS Records timer

[Timer]
OnStartupSec=5m
Unit=gad@%i.service

[Install]
WantedBy=timers.target