# roles/gad/templates/gad@.service.j2
[Unit]
Description=Gandi automatic DNS record updater
After=network-pre.target

[Service]
Type=simple
User=gandi
Group=gandi
RuntimeMaxSec=120s

[Install]
WantedBy=multi-user.target